# examples/systemd/legible.service
[Unit]
Description=Legible Sync Daemon
Documentation=https://github.com/platinummonkey/legible
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i

# Path to the legible binary
# Update this to match your installation location
ExecStart=/usr/local/bin/legible daemon --interval 10m --health-addr :8080

# Configuration file (optional)
# ExecStart=/usr/local/bin/legible --config /home/%i/.legible.yaml daemon

# Working directory
WorkingDirectory=/home/%i

# Environment variables (optional)
# Environment="LEGIBLE_OUTPUT_DIR=/home/%i/Documents/Legible"
# Environment="LEGIBLE_LOG_LEVEL=info"

# Restart policy
Restart=on-failure
RestartSec=30s

# Security settings (recommended)
# Restrict what the service can access
PrivateTmp=true
NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=read-only
# ReadWritePaths=/home/%i/Documents/Legible /home/%i/.legible

# Resource limits (optional)
# Limit CPU usage (100% = 1 core, 200% = 2 cores)
# CPUQuota=100%
# Limit memory usage
# MemoryLimit=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=legible

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
