# examples/systemd/remarkable-sync.service
[Unit]
Description=reMarkable Sync Daemon
Documentation=https://github.com/platinummonkey/remarkable-sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i

# Path to the remarkable-sync binary
# Update this to match your installation location
ExecStart=/usr/local/bin/remarkable-sync daemon --interval 10m --health-addr :8080

# Configuration file (optional)
# ExecStart=/usr/local/bin/remarkable-sync --config /home/%i/.remarkable-sync.yaml daemon

# Working directory
WorkingDirectory=/home/%i

# Environment variables (optional)
# Environment="REMARKABLE_SYNC_OUTPUT_DIR=/home/%i/Documents/remarkable"
# Environment="REMARKABLE_SYNC_LOG_LEVEL=info"

# Restart policy
Restart=on-failure
RestartSec=30s

# Security settings (recommended)
# Restrict what the service can access
PrivateTmp=true
NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=read-only
# ReadWritePaths=/home/%i/Documents/remarkable /home/%i/.remarkable-sync

# Resource limits (optional)
# Limit CPU usage (100% = 1 core, 200% = 2 cores)
# CPUQuota=100%
# Limit memory usage
# MemoryLimit=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remarkable-sync

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
